# packages/ftp-servers/proftpd/files/systemd/proftpd.service
[Unit]
Description=ProFTPd FTP Server
After=syslog.target network.target

[Service]
EnvironmentFile=/etc/conf.d/proftpd
ExecStart=/usr/sbin/proftpd -n $PROFTPD_OPTIONS

[Install]
WantedBy=multi-user.target
