# packages/net-analyzer/zabbix-agent/files/systemd/zabbix-agentd.service
# Copyright 2012 Kevin Decherf <kevin@kdecherf.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=Zabbix Agent
After=network.target

[Service]
ExecStart=/usr/sbin/zabbix_agentd
Type=forking
PIDFile=/run/zabbix/zabbix-agentd.pid

[Install]
WantedBy=multi-user.target
