# packages/sys-cluster/glance/files/systemd/glance-api.service
[Unit]
Description=Glance API service
After=network.target keystone.service

[Service]
ExecStart=/usr/bin/glance-api

[Install]
WantedBy=glance.target
