# packages/sys-cluster/glance/files/systemd/glance-registry.service
# Copyright 2012 Marc-Antoine Perennou <Marc-Antoine@Perennou.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=Glance registry service
After=network.target keystone.service mysql.service postgresql.service

[Service]
ExecStart=/usr/bin/glance-registry

[Install]
WantedBy=glance.target
