# packages/sys-cluster/glance/files/systemd/glance.target
[Unit]
Description=Glance System
Requires=network.target keystone.service postgresql.servive
After=network.target keystone.service postgresql.service
RefuseManualStart=no

[Install]
WantedBy=multi-user.target
