# packages/sys-cluster/glance/files/systemd/glance-registry.service
[Unit]
Description=Glance registry service
Requires=keystone.service
After=network.target keystone.service mysql.service postgresql.service

[Service]
ExecStart=/usr/bin/glance-registry

[Install]
WantedBy=glance.target
