# packages/ftp-servers/pure-ftpd/files/systemd/pure-ftpd.service
[Unit]
Description=Pure-FTPd Per-Connection
After=syslog.target

[Service]
EnvironmentFile=/etc/conf.d/pure-ftpd.conf
ExecStart=/usr/sbin/pure-ftpd $LISTEN $AUTH $OPTS

[Install]
WantedBy=multi-user.target
