# packages/ftp-servers/pure-ftpd/files/systemd/pure-ftpd.socket
[Unit]
Conflicts=pure-ftpd.service

[Socket]
ListenStream=21
Accept=yes

[Install]
WantedBy=sockets.target
