# packages/ftp-servers/pure-ftpd/files/systemd/pure-ftpd@.service
[Unit]
Description=Pure-FTPd Service
After=syslog.target

[Service]
EnvironmentFile=/etc/conf.d/pure-ftpd.conf
ExecStart=/usr/sbin/pure-ftpd $AUTH $OPTS
StandardInput=socket
