# packages/dev-db/etcd/files/systemd/etcd.service.d/timers.conf
[Install]
Also=etcd-compact.timer etcd-defrag.timer etcd-snapshot.timer
