# packages/net-dns/pdns/files/systemd/pdns.service
# Copyright 2012 Kevin Decherf <kevin@kdecherf.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=PowerDNS Authoritative Server
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/sbin/pdns_server
ExecReload=/usr/bin/pdns_control reload
ExecStop=/usr/bin/pdns_control quit
Restart=on-failure
RestartSec=2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
