# packages/net/rabbitmq-server/files/systemd/rabbitmq-server.service
[Unit]
Description=RabbitMQ messaging service
After=network.target

[Service]
User=rabbitmq
ExecStart=/usr/sbin/rabbitmq-server

[Install]
WantedBy=multi-user.target
