# packages/net/rabbitmq-server/files/systemd/rabbitmq-server.service
# Copyright 2012 Marc-Antoine Perennou <Marc-Antoine@Perennou.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=RabbitMQ messaging service
After=network.target

[Service]
User=rabbitmq
WorkingDirectory=/var/lib/rabbitmq
ExecStart=/usr/sbin/rabbitmq-server

[Install]
WantedBy=multi-user.target
