# packages/net-analyzer/zabbix-agent2/files/systemd/zabbix-agent2.service
[Unit]
Description=Zabbix Agent 2
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/zabbix_agent2 -f -c /etc/zabbix/zabbix_agent2.conf
RuntimeDirectory=zabbix
PIDFile=/run/zabbix/zabbix_agent2.pid
User=zabbix

[Install]
WantedBy=multi-user.target
