# packages/net/sozu-prometheus-connector/files/systemd/sozu-prometheus-connector-0.15.service
[Unit]
Description = Sōzu prometheus connector - Exposes Sōzu metrics in a prometheus format on HTTP.
Documentation = https://github.com/CleverCloud/sozu-prometheus-connector
After=network-online.target local-fs.target remote-fs.target time-sync.target
Wants=network-online.target local-fs.target remote-fs.target time-sync.target

StartLimitIntervalSec=30
StartLimitBurst=64

[Service]
Type=simple

ExecStart=/usr/libexec/sozu-prometheus-connector.0.15 --config /etc/sozu/0.15/connector/prometheus.toml
Restart=always

TimeoutStopSec=30
LimitNOFILE=65536

# Since systemd 229, should be in [Unit] but in order to support systemd <229,
# it is also supported to have it here.
StartLimitInterval=30
StartLimitBurst=64

[Install]
WantedBy=sozu-0.15.service
