# packages/dev-db/pganalyze-collector/files/systemd/pganalyze-collector.service
[Unit]
Description=pganalyze-collector
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStart="/usr/bin/pganalyze-collector"
Restart=always

[Install]
WantedBy=multi-user.target
