# packages/net-dns/pdns-recursor/files/systemd/pdns-recursor.service
# Copyright 2012 Kevin Decherf <kevin@kdecherf.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=PowerDNS Recursive Server
Documentation=man:pdns_recursor(1) man:rec_control(1)
Wants=network-online.target nss-lookup.target
Before=nss-lookup.target
After=network-online.target

[Service]
ExecStart=/usr/sbin/pdns_recursor --daemon=no
PrivateTmp=true

[Install]
WantedBy=multi-user.target
