# packages/sys-cluster/ceph/files/systemd/ceph-mds@.service
# Copyright 2012 Kevin Decherf <kevin@kdecherf.com>
# Distributed under the termes of the GNU Geneal Public License v2
#     Based upon ArchLinux Ceph-0.48.2 systemd service

[Unit]
Description=Ceph mds.%i
Requires=network.target
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=2s
ExecStart=/usr/bin/ceph-mds -f -i %i

[Install]
WantedBy=multi-user.target
