# packages/dev-db/etcd/files/systemd/etcd-snapshot.timer
[Unit]
Description = Trigger snapshot of etcd's key space
Documentation = https://etcd.io/docs/v3.5/op-guide/maintenance/

[Timer]
OnCalendar = 00:00:00
Persistent = true

[Install]
WantedBy = timers.target
