# packages/dev-db/etcd/files/systemd/etcd-compact.timer
[Unit]
Description = Trigger compaction of etcd's key space
Documentation = https://etcd.io/docs/v3.5/op-guide/maintenance/

[Timer]
OnCalendar = 00/3:00:00
Persistent = true

[Install]
WantedBy = timers.target
