# packages/dev-db/Percona-Server/files/systemd/mysql.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=MySQL database server (single instance)
After=syslog.target
After=network.target

# The bl-1 init script contains a mass of possible configurations and options to
# run several MySQL server instances, define nice-ness and configure everything
# up to different defaults for databases on an SSD.
# In most usecases, a single server instance is enough. That's what this service
# does; it starts a single server instance with sensible defaults.
# If you need more, provide a patch.

[Service]
PrivateTmp=true
User=mysql
Group=mysql

ExecStart=/usr/sbin/mysqld --defaults-file=/etc/mysql/my.cnf --basedir=/usr --datadir=/var/lib/mysql --pid-file=/run/mysqld/mysqld.pid --socket=/run/mysqld/mysqld.sock --log-error=/var/log/mysql/mysql.log -W

# Give a reasonable amount of time for the server to start up/shut down
TimeoutSec=300

# We rely on systemd, not mysqld_safe, to restart mysqld if it dies
Restart=always

[Install]
WantedBy=multi-user.target
