# packages/dev-db/pgbouncer/files/systemd/pgbouncer.service
[Unit]
Description=PgBouncer
After=network.target

[Service]
User=postgres
ExecStart=/usr/bin/pgbouncer /etc/pgbouncer/pgbouncer.ini
ExecReload=/usr/bin/pgbouncer -R /etc/pgbouncer/pgbouncer.ini

[Install]
WantedBy=multi-user.target
