# packages/dev-db/redis/files/systemd/redis.service
# Copyright 2012 Marc-Antoine Perennou <Marc-Antoine@Perennou.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=Redis Server
After=network.target
Requires=network.target

[Service]
User=redis
WorkingDirectory=/var/db/redis
ExecStart=/usr/bin/redis-server /etc/redis.conf
Restart=always

[Install]
WantedBy=multi-user.target
