# packages/dev-db/etcd/files/etcd-defrag.service
[Unit]
Description = Snapshot etcd's key-space
Documentation = https://etcd.io/docs/v3.5/op-guide/maintenance/
After = etcd.service etcd-compact.service
Wants = etcd.service

[Service]
User = etcd
Type = oneshot
ExecStart = /usr/bin/etcdctl defrag
Restart = on-failure
RestartSec = 60s

[Install]
WantedBy = multi-user.target
