# packages/net/memcached/files/systemd/memcached.service
# Copyright 2012 Kevin Decherf <kevin@kdecherf.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=Memcached
After=network.target

[Service]
ExecStart=/usr/bin/memcached -umemcached

[Install]
WantedBy=multi-user.target
