# packages/sys-auth/keystone/files/systemd/keystone.service
# Copyright 2012 Marc-Antoine Perennou <Marc-Antoine@Perennou.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=Keystone service
After=network.target mysql.service postgresql.service

[Service]
ExecStart=/usr/bin/keystone-all

[Install]
WantedBy=multi-user.target
