# packages/sys-cluster/nova/files/systemd/nova-compute.service
# Copyright 2012 Marc-Antoine Perennou <Marc-Antoine@Perennou.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=Start nova-compute node
Requires=network.target libvirtd.service iptables.service
After=network.target keystone.service mysql.service postgresql.service glance.target rabbitmq-server.servive libvirtd.service iptables.service

[Service]
Type=simple
ExecStart=/usr/bin/nova-compute

[Install]
WantedBy=nova.target
