# packages/sys-cluster/nova/files/systemd/nova-api.service
[Unit]
Description=Start nova-api node
Requires=network.target keystone.service postgresql.servive glance.target rabbitmq-server.servive
After=network.target keystone.service postgresql.service glance.target rabbitmq-server.servive

[Service]
Type=simple
ExecStart=/usr/bin/nova-api

[Install]
WantedBy=nova.target
