# packages/sys-cluster/nova/files/systemd/nova-objectstore.service
[Unit]
Description=Start nova-objectstore node
Requires=network.target
After=network.target keystone.service mysql.service postgresql.service glance.target rabbitmq-server.servive

[Service]
Type=simple
ExecStart=/usr/bin/nova-objectstore

[Install]
WantedBy=nova.target
