# packages/sys-cluster/nova/files/systemd/nova-cert.service
# Copyright 2012 Marc-Antoine Perennou <Marc-Antoine@Perennou.com>
# Distributed under the termes of the GNU Geneal Public License v2

[Unit]
Description=Start nova-cert node
Requires=network.target
After=network.target keystone.service mysql.service postgresql.service glance.target rabbitmq-server.servive

[Service]
Type=simple
ExecStart=/usr/bin/nova-cert

[Install]
WantedBy=nova.target
