# packages/sys-cluster/nova/files/systemd/nova-scheduler.service
[Unit]
Description=Start nova-scheduler node
Requires=network.target
After=network.target keystone.service mysql.service postgresql.service glance.target rabbitmq-server.servive

[Service]
Type=simple
ExecStart=/usr/bin/nova-scheduler

[Install]
WantedBy=nova.target
