# packages/sys-cluster/nova/files/systemd/nova-compute.service
[Unit]
Description=Start nova-compute node
Requires=network.target
After=network.target keystone.service mysql.service postgresql.service glance.target rabbitmq-server.servive libvirtd.service

[Service]
Type=simple
ExecStart=/usr/bin/nova-compute

[Install]
WantedBy=nova.target
