# packages/sys-cluster/nova/files/systemd/nova-network.service
[Unit]
Description=Start nova-network node
Requires=network.target
After=network.target keystone.service mysql.service postgresql.service glance.target rabbitmq-server.servive

[Service]
Type=simple
ExecStart=/usr/bin/nova-network

[Install]
WantedBy=nova.target
